# ssh-bouncer.service
[Unit]
Description=ssh-bouncer service
WantedBy=multi-user.target
 
[Service]
Type=simple
ExecStart=/usr/local/bin/ssh-bouncer
Restart=always